# disk/src/disk-widget.timer
[Unit]
Description=Planifie la collecte du widget disque (toutes les 15 min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
